# webex/single_sign_off/challenge/supervisord.conf

[supervisord]
logfile=stdout
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
stderr_logfile=stdout
stdout_logfile=stdout
user=root

[program:services]
command=/app/start.sh
autostart=true
autorestart=true
stderr_logfile=stdout
stdout_logfile=stdout
user=root
